# deploy/systemd/license-renewals.timer
# deploy/systemd/license-renewals.timer
[Unit]
Description=Daily license renewal pass

[Timer]
OnCalendar=*-*-* 03:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
